# phosphor-virtual-sensor.service.in
[Unit]
Description=Virtual sensors
Wants=obmc-mapper.target
After=obmc-mapper.target

[Service]
ExecStart=@bindir@/virtual-sensor
Type=dbus
BusName=xyz.openbmc_project.VirtualSensor
SyslogIdentifier=phosphor-virtual-sensor

[Install]
WantedBy=multi-user.target
